# src/deb/control/postrm
#!/bin/sh
# postrm script for #PACKAGE#
#
# see: dh_installdeb(1)

set -e

# summary of how this script can be called:
#        * <postrm> `remove'
#        * <postrm> `purge'
#        * <old-postrm> `upgrade' <new-version>
#        * <new-postrm> `failed-upgrade' <old-version>
#        * <new-postrm> `abort-install'
#        * <new-postrm> `abort-install' <old-version>
#        * <new-postrm> `abort-upgrade' <old-version>
#        * <disappearer's-postrm> `disappear' <overwriter>
#          <overwriter-version>
# for details, see http://www.debian.org/doc/debian-policy/ or
# the debian-policy package


get_hostname() {
    HOSTNAME=`hostname --all-fqdns`
    if [ -z ${HOSTNAME}  ]
    then
        HOSTNAME=`hostname`
    fi

    if [ -z ${HOSTNAME}  ] 
    then
        HOSTNAME="localhost"
    fi

    HOSTNAME="$(echo "${HOSTNAME}" | tr -d '[[:space:]]')"    
    # return hostname
    echo "$HOSTNAME"
}


case "$1" in

    remove|abort-install|disappear)
  
    
    ;;
    purge)
    # clean ldap
    ODN_CONF_LDAP=/etc/ldap/slapd.d/cn\=config/olcDatabase\=\{2\}hdb.ldif
    if [ -f $ODN_CONF_LDAP  ] ; then     
        rm -f $ODN_CONF_LDAP 
    fi
    
    ODN_CONF_LDAP=/etc/ldap/slapd.d/cn\=config/olcDatabase\=\{2\}hdb
    if [ -d $ODN_CONF_LDAP  ] ; then     
        rm -rf $ODN_CONF_LDAP 
    fi
  
    LDAP_DB_PATH=/var/lib/ldap_odn
    if [ -d $LDAP_DB_PATH  ] ; then     
        rm -rf $LDAP_DB_PATH 
    fi
    
    # reload default slapd configuration
    if [ -x "/etc/init.d/slapd" ]; then
        invoke-rc.d slapd restart || true
    fi 
    
    # clean cert
    HOSTNAME=$(get_hostname)
    CERTS_PATH=/usr/local/share/ca-certificates/$HOSTNAME
    if [ -d $CERTS_PATH  ] ; then     
        rm -rf $CERTS_PATH 
    fi
    
    update-ca-certificates

    ;;
    
    
    upgrade|failed-upgrade|abort-upgrade)

    ;;

    *)
        echo "postrm called with unknown argument \`$1'" >&2
        exit 1
    ;;
esac

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

#DEBHELPER#

exit 0
