# src/deb/control/postinst
#!/bin/sh
# postinst script for #PACKAGE#
#
# see: dh_installdeb(1)

set -e

# summary of how this script can be called:
#        * <postinst> `configure' <most-recently-configured-version>
#        * <old-postinst> `abort-upgrade' <new version>
#        * <conflictor's-postinst> `abort-remove' `in-favour' <package>
#          <new-version>
#        * <postinst> `abort-remove'
#        * <deconfigured's-postinst> `abort-deconfigure' `in-favour'
#          <failed-install-package> <version> `removing'
#          <conflicting-package> <version>
# for details, see http://www.debian.org/doc/debian-policy/ or
# the debian-policy package


create_user() {
  if ! id odn-cas > /dev/null 2>&1 ; then
    adduser --system --no-create-home --disabled-password --home /usr/share/odn-cas odn-cas --shell /bin/false
  fi
}

fix_perms() {
    chmod +x  /etc/init.d/odn-cas
    chown  odn-cas /var/cache/odn-cas
    chown  odn-cas /var/tmp/odn-cas
    chown  odn-cas /var/log/odn-cas
}

service_start() {
    update-rc.d odn-cas defaults

    if [ -x "`which invoke-rc.d 2>/dev/null`" ]; then
        invoke-rc.d odn-cas restart || true
    else
        service odn-cas restart || true
    fi
}

set_ldap() {
    echo ">> set_ldap"
    LDAP_DB_PATH=/var/lib/ldap_odn
    chown -R openldap:openldap $LDAP_DB_PATH
    chown openldap:openldap /etc/odn-simple/ldap/olcDatabase\=\{2}hdb.ldif
    #
    if [ ! -e  /etc/ldap/slapd.d/cn\=config/olcDatabase\=\{2\}hdb.ldif  ]; then
        echo "set ldap configuration"
        ln -s /etc/odn-simple/ldap/olcDatabase\=\{2}hdb.ldif  /etc/ldap/slapd.d/cn\=config/olcDatabase\=\{2\}hdb.ldif
    fi
        
    invoke-rc.d slapd restart
    TREE="dc=opendata,dc=org"
    echo "ldap tree: ${TREE} for IDM created"
    ldapadd -x -D cn=admin,$TREE -w 5ecr3t -f  /usr/share/odn-simple/ldap/admin.ldif || true
    ldapadd -x -D cn=admin,$TREE -w 5ecr3t -f  /usr/share/odn-simple/ldap/people.ldif || true
    echo "ldap user for IDM created"
    echo "<< set_ldap"
} 


case "$1" in
    configure)
        create_user
        fix_perms
        set_ldap
        service_start
    ;;

    abort-upgrade|abort-remove|abort-deconfigure)
    ;;

    *)
        echo "postinst called with unknown argument \`$1'" >&2
        exit 1
    ;;
esac

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

#DEBHELPER#

exit 0
